# puppet/files/ceilometer/init/ceilometer-api.conf
description "Ceilometer API Server"
author "Longgeek <longgeek@gmail.com>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -s /bin/sh -c "exec ceilometer-api --config-file=/etc/ceilometer/ceilometer.conf" root
