# puppet/files/trove/init/trove-api.conf
description "Trove API Server"
author "Longgeek <longgeek@gmail.com>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -s /bin/sh -c "exec /livestack/bin/trove-api --config-file=/etc/trove/trove.conf" root
