# puppet/files/swift/init/swift-object-server.conf
description "Swift Object Server"
author "Longgeek <longgeek@gmail.com>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -s /bin/sh -c "exec /livestack/bin/swift-object-server /etc/swift/object-server.conf" root
