# puppet/files/neutron/init/neutron-lbaas-agent.conf
description "Neutron LBAAS Agent Server"
author "Longgeek <longgeek@gmail.com>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -s /bin/sh -c "exec neutron-lbaas-agent --config-file /etc/neutron/neutron.conf --config-file=/etc/neutron/services/loadbalancer/haproxy/lbaas_agent.ini --log-file=/var/log/neutron/lbaas-agent.log" root
