# puppet/files/glance/init/glance-api.conf
description "Glance API Server"
author "Longgeek <longgeek@gmail.com>"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

respawn

exec su -s /bin/sh -c "exec /livestack/bin/glance-api --config-file=/etc/glance/glance-api.conf" root
